# etc/services.d/medusa/run
#!/usr/bin/with-contenv /bin/sh

exec \
        s6-setuidgid medusa ${PYTHON_INTERPRETER} ${APP_PATH}/start.py \
        --nolaunch --datadir /config
