# root/etc/s6-overlay/s6-rc.d/medusa/run
#! /usr/bin/with-contenv /bin/sh
# shellcheck shell=sh

exec s6-setuidgid medusa ${PYTHON_INTERPRETER} ${APP_PATH}/start.py \
        --nolaunch --datadir /config
